# templates/gunicorn/gunicorn.conf.erb
# gunicorn

description "gunicorn server"
author "Calen Pennington <cpennington@mitx.mit.edu>"

<% if scope.lookupvar('stacked') == false -%>
start on runlevel [2345]
stop on runlevel [!2345]
<% else -%>
start on started edxapp 
stop on stopped edxapp 
<% end -%>

respawn
<% if scope.lookupvar('respawn_limit') != false -%>
respawn limit 3 30
<% end -%>

env PID=/var/tmp/<%= name %>.pid
env NEW_RELIC_CONFIG_FILE=<%= scope.lookupvar('base') %>/newrelic.ini
env NEWRELIC=<%= scope.lookupvar('virtualenv') %>/bin/newrelic-admin
<% if scope.lookupvar('workers') -%>
env WORKERS=<%= scope.lookupvar('workers') %>
<% else -%>
env WORKERS=<%= 4 * scope.lookupvar('::processorcount').to_i %>
<% end -%>
env PORT=<%= scope.lookupvar('port') %>
env LANG=en_US.UTF-8
env DJANGO_SETTINGS_MODULE=<%= scope.lookupvar('settings_module') %>
env SERVICE_VARIANT="<%= scope.lookupvar('title') %>"

<% if scope.lookupvar('stacked') == false -%>
pre-start script
  find ${package_root} -user ${user} -type f -name '*.pyc' -delete || true
  find /tmp -user ${user} -type d -name tmp*mako -exec rm -rf {} \\; || true
end script
<% end -%>

chdir <%= scope.lookupvar('package_root') %>
<% if scope.lookupvar('app_interface') != 'python' -%>
setuid <%= scope.lookupvar('user') %>
<% end -%>

<% case scope.lookupvar('app_interface') when 'django' -%>
exec <% if scope.lookupvar('edx::newrelic::reporting') == true -%>$NEWRELIC run-program <% end -%><%= scope.lookupvar('virtualenv') %>/bin/gunicorn_django -b 127.0.0.1:$PORT -w $WORKERS --timeout=<%= scope.lookupvar('timeout') %> --pythonpath=<%= scope.lookupvar('package_root') %> --settings=<%= scope.lookupvar('settings_module') %>
<% when 'wsgi' -%>
exec <% if scope.lookupvar('edx::newrelic::reporting') == true -%>$NEWRELIC run-program <% end -%><%= scope.lookupvar('virtualenv') %>/bin/gunicorn --preload -b 127.0.0.1:$PORT -w $WORKERS --timeout=<%= scope.lookupvar('timeout') %> --pythonpath=<%= scope.lookupvar('package_root') %> <%= scope.lookupvar('wsgi_app') %>
<% when 'python' -%>
exec su - www-data -c '<%= scope.lookupvar('virtualenv') %>/bin/python <%= scope.lookupvar('package_root') %>/<%= scope.lookupvar('script_name') %> >/dev/null 1>&1'
<% end -%>

<% if (scope.lookupvar('title') == 'lms-xml') || (scope.lookupvar('title') == 'lms') || (scope.lookupvar('title') == 'lms-preview') -%>
post-start script
  while true
  do
    if $(curl -s -i localhost:$PORT/heartbeat | egrep -q '200 OK'); then
      break;
    else
      sleep 1;
    fi
  done
end script
<% end -%>